# data/entrance.service.in
[Unit]
Description=Enlightenment Display Manager
After=systemd-user-sessions.service

[Service]

Environment="LD_LIBRARY_PATH=@libdir@" "PATH=$PATH:@bindir@"

ExecStart=@SBINDIR@/entrance --nodaemon
Restart= always

[Install]
Alias=display-manager.service
