# data/entrance.service.in
[Unit]
Description=Enlightenment Display Manager
Conflicts=getty@tty1.service plymouth-quit.service
After=systemd-user-sessions.service getty@tty1.service plymouth-quit.service

[Service]

Environment="LD_LIBRARY_PATH=@libdir@"
IgnoreSIGPIPE=no
ExecStart=@SBINDIR@/entrance --nodaemon
Restart= always
StandardOutput=syslog
StandardError=inherit

[Install]
Alias=display-manager.service
